# b2evolution.cron
* * * * * root /usr/bin/php /opt/b2evolution/cron/cron_exec.php
